# etc/systemd/system/timers.target.wants/snap.certbot.renew.timer
[Unit]
# Auto-generated, DO NOT EDIT
Description=Timer renew for snap application certbot.renew
Requires=snap-certbot-2772.mount
After=snap-certbot-2772.mount
X-Snappy=yes

[Timer]
Unit=snap.certbot.renew.service
OnCalendar=*-*-* 04:04
OnCalendar=*-*-* 14:08

[Install]
WantedBy=timers.target
